# crontab
0 12 * * * root cd /usr/src/app && ./kindle.sh
